# etc/crontab
00 04 * * * /cg_tools/backup.sh >> /dev/stdout 2>&1
# crontab requires an empty line at the end of the file
